# logstash5/docker/service/logstash/run
#!/bin/sh

exec /sbin/setuser app /usr/share/logstash/bin/logstash --path.settings /etc/logstash