# dockbeat-client/docker/service/dockbeat/run
#!/bin/bash
exec dockbeat -c /etc/dockbeat/dockbeat.yml
