# monitoring-proxy/docker/service/redis/run
#!/bin/sh

exec /sbin/setuser redis redis-server /etc/redis/redis.conf
