# sensu-api/docker/service/sensu-api/run
#!/bin/sh

exec /opt/sensu/bin/sensu-api -vc /etc/sensu/config.json -d /conf.d
