# sensu-client/docker/service/sensu-client/run
#!/bin/sh

exec /opt/sensu/bin/sensu-client -vc /etc/sensu/config.json -d /conf.d
