# monitoring-proxy/docker/service/bucky/run
#!/bin/sh

exec /sbin/setuser bucky bucky /etc/bucky/bucky.conf
