# monitoring-proxy/docker/service/logstash/run
#!/bin/sh

exec /sbin/setuser logstash /opt/logstash/bin/logstash -f /etc/logstash/logstash.conf
